# scripts/external/systemd/runmap-tiles.timer.tmpl
[Unit]
Description=Schedule RunMap tiles builder

[Timer]
OnCalendar=hourly
Persistent=true
Unit=runmap-tiles.service

[Install]
WantedBy=timers.target
